# templates/graphite.service
[Unit]
Description=graphite container
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker start -a graphite
ExecStop=/usr/bin/docker stop -t 2 graphite

[Install]
WantedBy=multi-user.target
